# deploy/public_proxy/systemd/bambuddy-cache-warmer.timer
[Unit]
Description=Refresh Bambuddy API cache and snapshots every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target

