# board/systemd/bambuddy-reverse-tunnel.service
[Unit]
Description=Bambuddy reverse SSH tunnel to public proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Edit these values before enabling the unit on your own board.
Environment=SSH_KEY_PATH=/root/.ssh/bambuddy_reverse_tunnel
Environment=TUNNEL_USER=bambuddy-tunnel
Environment=PUBLIC_SERVER_HOST=public-server.example.com
Environment=QUEUE_SERVER_PORT=18001
Environment=PUBLIC_SERVER_BOARD_PORT=18000
Environment=BOARD_APP_PORT=8000
ExecStart=/usr/bin/ssh -NT \
  -o BatchMode=yes \
  -o ExitOnForwardFailure=yes \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=yes \
  -i ${SSH_KEY_PATH} \
  -L 127.0.0.1:${QUEUE_SERVER_PORT}:127.0.0.1:${QUEUE_SERVER_PORT} \
  -R 127.0.0.1:${PUBLIC_SERVER_BOARD_PORT}:127.0.0.1:${BOARD_APP_PORT} \
  ${TUNNEL_USER}@${PUBLIC_SERVER_HOST}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
